# etc/sysbackup-amp.cron.daily
#!/bin/sh
exec sysbackup-amp --quiet --syslog --auto --multithreaded --folder=/var/backups/sysbackup-amp
